# deploy/supervisor/webeye.conf
[program:webeye]
command=/usr/local/bin/gunicorn -w10 -b127.0.0.1:5000 app:app
directory=/webeye
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true